# CascadeProjects/windsurf-project/com.supermegabot.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.supermegabot.launcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rudolfsarkany/supermegabot-windsurf-agents/CascadeProjects/windsurf-project/auto-start.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rudolfsarkany/supermegabot-windsurf-agents/CascadeProjects/windsurf-project</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/supermegabot-launcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/supermegabot-launcher-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
