# CascadeProjects/windsurf-project/com.supermegabot.watchdog-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.supermegabot.watchdog-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/rudolfsarkany/supermegabot-windsurf-agents/CascadeProjects/windsurf-project/watchdog-monitor-simple.cjs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rudolfsarkany/supermegabot-windsurf-agents/CascadeProjects/windsurf-project</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/tmp/supermegabot-watchdog-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/supermegabot-watchdog-monitor-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
